# extras/cron.d/mutalyzer-ucsc-update
25 6 * * * www-data <MUTALYZER_BIN_UCSC_UPDATE>
